# scripts/05_automation/systemd/asdf.service
[Unit]
Description=Service de Ejercicio 5

[Service]
Type=Oneshot
ExecStart=/usr/bin/asdf

