# scripts/05_automation/systemd/asdf.service
[Unit]
Description=Ejecuta /usr/bin/asdf

[Service]
Type=oneshot
ExecStart=/usr/bin/asdf
User=ignacio
Group=ignacio

