# scripts/05_automation/systemd/backup.service
[Unit]
Description=Copia de seguridad entre semana

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup.sh
User=root
Group=root

