# scripts/05_automation/systemd/backup.service
[Unit]
Description=Gestion de usuarios desde CSV

[Service]
Type=oneshot
ExecStart=../backup.sh

