# scripts/05_automation/systemd/backup.service
[Unit]
Description=Ejecucion del backup

[Service]
Type=oneshot
ExecStart=/root/ASO_MDC/scripts/05_automation
