# scripts/05_automation/systemd/gestion_usuarios.service
[Unit]
Description=Copia de seguridad entre semana

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gestion_usuarios.sh
User=root
Group=root

