# scripts/05_automation/systemd/gestion_usuarios.service
[Unit]
Description=Gestion automática de usuarios desde CSV

[Service]
Type=oneshot
ExecStart=/root/gestion_usuarios.sh
User=root
