# scripts/05_automation/systemd/gestion_usuarios.timer
[Unit]
Description=Timer para gestion_usuarios

[Timer]
OnCalendar=Mon..Fri 07:00
OnCalendar=Mon..Fri 14:30

[Install]
WantedBy=timers.target

