# scripts/05_automation/systemd/script.service
[Unit]
Description=Script para añadir o borrar usuario

[Service]
Type=oneshot
ExecStart=/usr/local/bin/script.sh
User=root
Group=root
