# scripts/05_automation/systemd/systemd_ejer5.service
[Unit]
Description=systemd para ejercicio 5

[Service]
Type=oneshot
ExecStart=/usr/bin/asdf
User=root
Group=root
