# scripts/05_automation/crontab
0 7 * * 1-5 backup.sh
30 14 * * 1-5 backup.sh
