# scripts/05_automation/crontab
0 7 * * 1-5 /root/gestion_usuarios.sh
30 14 * * 1-5 /root/gestion_usuarios.sh
