# scripts/05_automation/gestion_usuarios.crontab
0 7 * * 6,0 root /home/hugo/Escritorio/ASO_ACT/usuarios.sh
30 14 * * 6,0 root /home/hugo/Escritorio/ASO_ACT/usuarios.sh
